# System/Library/LaunchDaemons/com.apple.SafariCloudHistoryPushAgent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.apple.SafariCloudHistoryPushAgent</string>
	<key>Program</key>
	<string>/usr/libexec/SafariCloudHistoryPushAgent</string>
	<key>UserName</key>
	<string>mobile</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.SafariCloudHistoryPushAgent</key>
		<true/>
		<key>com.apple.SafariCloudHistoryPushAgent.Push</key>
		<true/>
	</dict>
	<key>EnablePressuredExit</key>
	<true/>
</dict>
</plist>
