# System/Library/LaunchDaemons/com.apple.mobile.installd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>GroupName</key>
	<string>mobile</string>
	<key>UserName</key>
	<string>mobile</string>
	<key>Label</key>
	<string>com.apple.mobile.installd</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.mobile.installd</key>
		<true/>
	</dict>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/installd</string>
		<string>--idleExit</string>
	</array>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>EnableTransactions</key>
	<true/>
	<key>RunAtLoad</key>
	<true/>
	<key>JetsamProperties</key>
	<dict>
		<key>JetsamMemoryLimit</key>
		<integer>5000</integer>
	</dict>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.mobile.developer_image_mounted</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobile.developer_image_mounted</string>
			</dict>
		</dict>
	</dict>
</dict>
</plist>
