# System/Library/LaunchDaemons/com.apple.coreduetd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnableTransactions</key>
	<true/>
	<key>EnablePressuredExit</key>
	<true/>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>Label</key>
	<string>com.apple.coreduetd</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.coreduetd.pairedsync</key>
		<true/>
		<key>com.apple.corespotlight.receiver.coreduet</key>
		<true/>
		<key>com.apple.coreduetd</key>
		<true/>
		<key>com.apple.coreduetd.knowledgebase</key>
		<true/>
		<key>com.apple.coreduetd.people</key>
		<true/>
		<key>com.apple.coreduetd.knowledge</key>
		<true/>
		<key>com.apple.coreduetd.batterysaver</key>
		<true/>
		<key>com.apple.coreduetd.context</key>
		<true/>
		<key>com.apple.private.alloy.coreduet-idswake</key>
		<true/>
	</dict>
	<key>Program</key>
	<string>/usr/libexec/coreduetd</string>
	<key>UserName</key>
	<string>mobile</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.bluetooth.connection</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.bluetooth.connection</string>
			</dict>
			<key>com.apple.MediaRemote.nowPlayingApplicationIsPlayingDidChange</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.MediaRemote.nowPlayingApplicationIsPlayingDidChange</string>
			</dict>
			<key>com.apple.coremedia.carplayisconnected</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.coremedia.carplayisconnected</string>
			</dict>
			<key>com.apple.coreaudio.RoutingConfiguration</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.coreaudio.RoutingConfiguration</string>
			</dict>
			<key>com.apple.coreduetd.nearbydeviceschanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.coreduetd.nearbydeviceschanged</string>
			</dict>
			<key>com.apple.springboard.lockstate</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.springboard.lockstate</string>
			</dict>
			<key>com.apple.springboard.pluggedin</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.springboard.pluggedin</string>
			</dict>
			<key>com.apple.coreduet.idslaunchnotification</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.coreduet.idslaunchnotification</string>
			</dict>
			<key>com.apple.icloud.fmip.lostmode.enable</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.icloud.fmip.lostmode.enable</string>
			</dict>
			<key>com.apple.system.powersources.percent</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.system.powersources.percent</string>
			</dict>
			<key>com.apple.system.powersources.source</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.system.powersources.source</string>
			</dict>
			<key>com.apple.system.powermanagement.uservisiblepowerevent</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.system.powermanagement.uservisiblepowerevent</string>
			</dict>
			<key>com.apple.coreaudio.borealisTrigger</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.coreaudio.borealisTrigger</string>
			</dict>
			<key>SignificantTimeChangeNotification</key>
			<dict>
				<key>Notification</key>
				<string>SignificantTimeChangeNotification</string>
			</dict>
			<key>com.apple.powermanagement.restartpreventers</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.powermanagement.restartpreventers</string>
			</dict>
			<key>com.apple.powermanagement.idlesleeppreventers</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.powermanagement.idlesleeppreventers</string>
			</dict>
			<key>com.apple.powermanagement.systemsleeppreventers</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.powermanagement.systemsleeppreventers</string>
			</dict>
			<key>com.apple.system.thermalpressurelevel</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.system.thermalpressurelevel</string>
			</dict>
			<key>com.apple.iokit.hid.displayStatus</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.iokit.hid.displayStatus</string>
			</dict>
			<key>com.apple.symptoms.materialLinkQualityChange</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.symptoms.materialLinkQualityChange</string>
			</dict>
			<key>com.apple.system.config.network_change</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.system.config.network_change</string>
			</dict>
			<key>com.apple.locationd.vehicle.connected</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.locationd.vehicle.connected</string>
			</dict>
			<key>com.apple.locationd.vehicle.disconnected</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.locationd.vehicle.disconnected</string>
			</dict>
			<key>com.apple.GeoServices.navigation.started</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.GeoServices.navigation.started</string>
			</dict>
			<key>com.apple.GeoServices.navigation.stopped</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.GeoServices.navigation.stopped</string>
			</dict>
		</dict>
	</dict>
</dict>
</plist>
