# System/Library/NanoLaunchDaemons/com.apple.bulletindistributord.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Disabled</key>
	<true/>
	<key>Label</key>
	<string>com.apple.bulletindistributord</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.bulletindistributor.testservice</key>
		<true/>
		<key>com.apple.bulletindistributor.pairedsync</key>
		<true/>
		<key>com.apple.bulletindistributord.server</key>
		<true/>
		<key>com.apple.NotificationsUIPreferredDeviceManager.systemstateconnection</key>
		<true/>
	</dict>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/bulletindistributord</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>EnablePressuredExit</key>
	<false/>
	<key>CFBundleAllowMixedLocalizations</key>
	<true/>
</dict>
</plist>
