# System/Library/NanoLaunchDaemons/com.apple.pairedunlockd-phone.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.apple.pairedunlockd</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.paired-unlock</key>
		<true/>
		<key>com.apple.private.alloy.pairedunlock-idswake</key>
		<true/>
	</dict>
	<key>UserName</key>
	<string>mobile</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/PairedUnlock.framework/pairedunlockd</string>
	</array>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>EnablePressuredExit</key>
	<true/>
	<key>Disabled</key>
	<true/>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>KeybagStateChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mobile.keybagd.lock_status</string>
			</dict>
		</dict>
	</dict>
</dict>
</plist>
