# System/Library/NanoLaunchDaemons/com.apple.nanomailbootstrapd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.apple.nanomailbootstrapd</string>
	<key>MachServices</key>
	<dict>
        <key>com.apple.private.alloy.mail.sync.messages-idswake</key>
        <true/>
        <key>com.apple.private.alloy.mail.sync.content-idswake</key>
        <true/>
        <key>com.apple.private.alloy.mail.sync.protected-idswake</key>
        <true/>
        <key>com.apple.private.alloy.mail.sync.protected.content-idswake</key>
        <true/>
        <key>com.apple.private.alloy.mail.sync.accounts-idswake</key>
        <true/>
        <key>com.apple.private.alloy.mail.fetches-idswake</key>
        <true/>
    </dict>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/NanoMailKitServer.framework/nanomailbootstrapd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>Disabled</key>
	<true/>
</dict>
</plist>
