# System/Library/NanoLaunchDaemons/com.apple.nanomailbootstrapd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.apple.nanomailbootstrapd</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/NanoMailKitServer.framework/nanomailbootstrapd</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
	<key>ProcessType</key>
	<string>Adaptive</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.mail.AutoLaunchSettingsChanged</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.mail.AutoLaunchSettingsChanged</string>
			</dict>
			<key>com.apple.nanomail.idslaunchnotification</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.nanomail.idslaunchnotification</string>
			</dict>
		</dict>
	</dict>
    <key>Disabled</key>
    <true/>
</dict>
</plist>
