# System/Library/NanoLaunchDaemons/com.apple.nanomapscd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.apple.nanomapscd</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.mapssync.pairedsync</key>
		<true/>
		<key>com.apple.nanomaps.xpc.Maps</key>
		<true/>
		<key>com.apple.nanomaps.xpc.GeoServices</key>
		<true/>
		<key>com.apple.nanomaps.xpc.GeoServices.Navigation</key>
		<true/>
		<key>com.apple.nanomaps.xpc.MapKit</key>
		<true/>
		<key>com.apple.private.alloy.maps.proxy-idswake</key>
		<true/>
		<key>com.apple.private.alloy.maps.sync-idswake</key>
		<true/>
	</dict>
	<key>UserName</key>
	<string>mobile</string>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/MapsSupport.framework/nanomapscd</string>
	</array>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>EnablePressuredExit</key>
	<true/>
	<key>Disabled</key>
	<true/>
</dict>
</plist>
